# root/etc/uci-defaults/99-modem-sms-forward-AK68
#!/bin/sh

[ ! -x /etc/init.d/modem-sms-forward-AK68 ] || /etc/init.d/modem-sms-forward-AK68 enable
exit 0
